# supervisord.conf
[supervisord]
nodaemon=true

[program:web]
command=/bin/bash -c "cd /srv && poetry run gunicorn librarian.wsgi:application 0.0.0.0:8000 --capture-output --log-level=info --access-logfile - --error-logfile -"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

[program:worker]
command=/bin/bash -c "cd /srv && make run-worker"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
